# roles/kubernetes/files/etcd.service
[Unit]
Description=Run etcd on machine

[Service]
Type=forking
ExecStart=/usr/local/bin/etcd

[Install]
WantedBy=multi-user.target
